# elasticsearch/supervisor/conf.d/elasticsearch.conf
[program:elasticsearch]
command=/usr/share/elasticsearch/bin/elasticsearch --default.path.home=$ES_HOME
autostart=true
autorestart=true
