# kibana/supervisor/conf.d/kibana.conf
[program:kibana]
command=/usr/share/kibana/bin/kibana
autostart=true
autorestart=true

[program:nginx]
command=/usr/sbin/nginx
autostart=true
autorestart=true
